# container-files/etc/services.d/mysql/run
#!/bin/execlineb -P

/usr/bin/mysqld --user=mysql --console
